# examples/snapplicator.service
[Unit]
Description=Duplication of Snapper Snapshots

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /path/to/snapplicator/snapplicator.py
